# Milestone 4/SupervisorCtl/supervisord.ini
[supervisord]
nodaemon=true

[program:postgres]
command=/usr/bin/postgres_service.sh run
autostart=true
autorestart=true
stderr_logfile=/var/log/postgres_err.log
stdout_logfile=/var/log/postgres_out.log
stopsignal=QUIT

[program:tomcat]
command=/usr/local/tomcat/bin/catalina.sh run
stdout_logfile=/usr/local/tomcat/logs/catalina.out
stderr_logfile=/usr/local/tomcat/logs/catalina.out
environment=JAVA_HOME="/usr/lib/jvm/default-jvm",JAVA_BIN="/usr/lib/jvm/default-jvm/bin"
autorestart=true
startsecs=0
stopsignal=INT
stopasgroup=true
killasgroup=true